# packaging/app-flexshare.cron
15 4 * * * root /usr/clearos/apps/flexshare/deploy/set-file-permissions
